# system/systemd/teleprompter-mirror.service
[Unit]
Description=Teleprompter Mirror signaling server (localhost:8047)
PartOf=graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
# Optional per-machine overrides (port, bind address). The leading '-' makes
# the file optional, so the service starts fine with the built-in defaults.
EnvironmentFile=-%h/.config/teleprompter-mirror/config.env
ExecStart=/usr/bin/python3 __PROJECT_DIR__/app/mirror-server.py
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
