# system/systemd/teleprompter-wifi-rebind.service
[Unit]
Description=Re-probe MT7601U WiFi adapter if driver fails to bind

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/teleprompter-wifi-rebind.sh
