# teleprompter-wifi-rebind.service
[Unit]
Description=Re-probe MT7601U WiFi adapter if driver fails to bind

[Service]
Type=oneshot
ExecStart=__PROJECT_DIR__/wifi-rebind.sh
